# files/supervisor/wso2-esb.ini
[program:wso2-esb]
user = wso2
autostart = true
command = /opt/wso2/esb/bin/wso2server.sh
stdout_logfile = /var/log/wso2/esb/common.log
stderr_logfile = /var/log/wso2/esb/error.log
